# company-enrichment.service
[Unit]
Description=Company Enrichment Streamlit App
After=network.target

[Service]
User=root
WorkingDirectory=/root/company-enrichment-tool
ExecStart=/root/company-enrichment-tool/venv/bin/streamlit run app.py --server.port 8501 --server.address 0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
